# files/etc/profile.d/aliases.sh
#!/usr/bin/env sh

alias kc="kubectl"

alias la="la -la"
alias laf="la -laF"
